# debian/assistant-button.service
[Unit]
Description=Assistant button listener

[Service]
Type=simple
ExecStart=/usr/libexec/assistant-button
Restart=on-failure
